# crontab
* * * * * python -u /app/app/update_kuma.py