# scripts/pelican-reboot/pelican-reboot.service
[Unit]
Description=Cleanly stop all Pelican game servers and reboot the host

[Service]
Type=oneshot
EnvironmentFile=/etc/pelican-reboot.env
ExecStart=/usr/local/bin/pelican_reboot.py shutdown
